# network_stability_test/service/network_stability_test-slave.service.d/override.conf
[Service]
Environment="NODE_CONFIG_DIR=/opt/test/pc_test_tools/network_stability_test/config"